# rpi-image/weltempfaenger/rootfs-overlay/etc/profile.d/motd.sh
#!/bin/sh
set -eu
echo ''
echo '=*=*= Weltempfänger, Willkommen! =*=*='
echo ''
